# docker/supervisor.conf
[program:liquid-uwsgi]
command = /usr/local/bin/uwsgi --ini /liquid/uwsgi.ini

[program:nginx]
command = /usr/sbin/nginx
